# yumi-usb-update.service
[Unit]
Description=Yumi USB Offline Update Check
After=local-fs.target

[Service]
Type=oneshot
ExecStart=%h/yumi-config/yumi-usb-update-check.sh
TimeoutStartSec=180
StandardOutput=journal
StandardError=journal
RemainAfterExit=no

[Install]
WantedBy=default.target
